# deploy/auto-trader.service
[Unit]
Description=Auto-Trader Bot
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/auto-trader
ExecStart=/home/ubuntu/auto-trader/venv/bin/python main.py
Restart=always
RestartSec=10
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
